# phicomm-k2p/files/etc/crontabs/root
0 0 * * * wget -O files/etc/smartdns/anti-ad-smartdns.conf https://raw.githubusercontent.com/privacy-protection-tools/anti-AD/master/anti-ad-smartdns.conf
